# dot_config/systemd/user/sync.service
[Unit]
Description=Sync personal files with a remote called drive

[Service]
ExecStart=%h/.scripts/sync.sh

[Install]
WantedBy=default.target
