# user-services/.config/systemd/user/sync.service
[Unit]
Description=Sync personal files with Google Drive

[Service]
ExecStart=%h/.scripts/sync.sh

[Install]
WantedBy=default.target
